# com.user.slack-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.slack-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/mpovolny/Projects/slack-search &amp;&amp; /opt/homebrew/bin/uv run slack-search refresh --curl "$(cat .curl)" --no-files 2&gt;&amp;1 | while IFS= read -r line; do printf '[%s] %s\n' "$(date '+%Y-%m-%d %H:%M:%S')" "$line"; done</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/slack-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/slack-refresh.log</string>
</dict>
</plist>
